# overlay/system/usr/lib/systemd/system/wpa_suplicant.service.d/devicehacks.conf
[Service]
ExecStart=
ExecStart=/sbin/wpa_supplicant -u -s -O /run/wpa_supplicant -c /etc/wpa_supplicant/wpa_supplicant_overlay.conf -m /etc/wpa_supplicant/p2p_supplicant_overlay.conf
